# mail/qmail/pkg/qmail.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/rc"
daemon_flags=""

. /etc/rc.d/rc.subr

rc_cmd $1
