# mail/qmail/pkg/qmail.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/rc"
daemon_flags=""

. /etc/rc.d/rc.subr

pexp="qmail-send"
daemon_flags="${daemon_flags} &"

rc_cmd $1
